# ops/watchdogs/laravel/laravel-platform-watchdog.timer
[Unit]
Description=Run Laravel platform watchdog every 3 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=3min
AccuracySec=30s
Persistent=true
Unit=laravel-platform-watchdog.service

[Install]
WantedBy=timers.target
