# ops/watchdogs/laravel/laravel-platform-watchdog.service
[Unit]
Description=Self-heal Laravel API and RAGFlow bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/laravel-platform-watchdog.sh
